# nodewatch.service
[Unit]
Description=Nodewatch

[Service]
ExecStart=/usr/local/bin/nodewatchd/nodewatchd/nodewatchd runserver

